# cookbooks/fuelphp/templates/default/rc.local.erb
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

# remount /mnt/fuelphp for changing gid
mount | grep vboxsf | grep -q fuelphp
if [ $? eq "1" ]; then
  umount /mnt/fuelphp
  mount -t vboxsf -o uid=`id -u vagrant`,gid=`getent group apache | cut -d: -f3`,dmode=775,fmode=664 /mnt/fuelphp /mnt/fuelphp
fi

touch /var/lock/subsys/local
