# cookbooks/fuelphp/templates/default/remount-mnt-fuelphp.conf.erb
start on vagrant-mounted
 
exec /etc/init.d/remount-mnt-fuelphp

